# systemd/acpi-wakeup-fix.service
[Unit]
Description=Manage ACPI Wakeup Events for Power Optimization
After=multi-user.target
# This service prevents PCIe/USB devices from waking up the CPU unnecessarily

[Service]
Type=oneshot
ExecStart=/usr/local/bin/acpi-wakeup-fix.sh
# Ensures the service stays active after running the script once
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target